# darkstar/s6/plumber/finish
#!/bin/sh
echo "[S6] Plumber/Valve process exited with code $1, signal $2"
