# r-runtime/s6/plumber/run
#!/bin/sh
exec Rscript /app/valve_launcher.R
