# r-runtime/s6/plumber/run
#!/bin/sh
# Drop privileges from root (PID 1) to ruser for the R process (HIGHSEC §4.1)
exec s6-setuidgid ruser Rscript /app/valve_launcher.R
