# darkstar/s6/finngen-artifacts-init/up
/etc/s6-overlay/scripts/finngen-artifacts-init.sh
