# darkstar/s6/app-readable-init/up
/etc/s6-overlay/scripts/app-readable-init.sh
